# scripts/agi-id.service
[Unit]
Description=Aionima Local ID Service — OAuth credential broker
Documentation=https://github.com/Civicognita/aionima-local-id
After=network.target agi.service
Wants=agi.service

[Service]
Type=simple
User=%AIONIMA_USER%
Group=%AIONIMA_USER%
WorkingDirectory=/opt/agi-local-id
ExecStartPre=/usr/bin/npx drizzle-kit migrate
ExecStart=/usr/bin/node dist/index.js
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=agi-id

# Environment file for secrets (DATABASE_URL, ENCRYPTION_KEY, OAuth credentials).
# Installed by `scripts/install.sh` alongside this unit; AGI owns its lifecycle.
EnvironmentFile=/opt/agi-local-id/.env

# --- BEGIN CREDENTIALS ---
# --- END CREDENTIALS ---

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/opt/agi-local-id/dist /opt/agi-local-id/drizzle

[Install]
WantedBy=multi-user.target
